# debian/iitm-login-manager/usr/etc/systemd/user/iitm-login-manager.service
[Unit]
Description=IITM Login Manager System Tray
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/iitm-login-tray
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
